# cmd/platform/templates/streamed.container
[Unit]
Description=centos-streamed web server (Datastar app)
Wants=network-online.target
After=network-online.target

[Container]
Image=streamed.build
ContainerName=streamed
Network=proxy.network
Environment={{env "SERVER_NAME" .ServerName}}
Environment={{env "SERVER_OS" .ServerOS}}
Environment={{env "SERVER_KERNEL" .ServerKernel}}

[Service]
Restart=always
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
